# vision-roi.service
[Unit]
Description=vision roi
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/jetson/SMARTSense/ninenox/VisionROI

# เรียก app.py โหมด uvicorn
ExecStart=/home/jetson/SMARTSense/ninenox/venv/bin/python /home/jetson/SMARTSense/ninenox/VisionROI/app.py --port 12000 --use-uvicorn

# หยุดแบบนุ่ม: ยิง /_quit ให้แอป cleanup เองก่อน
ExecStop=/usr/bin/curl -fsS -X POST http://127.0.0.1:12000/_quit
ExecStop=/bin/kill -s SIGTERM $MAINPID

User=jetson
Group=jetson
Environment=PYTHONUNBUFFERED=1

# ช่วยให้ RTSP/FFmpeg ไม่ค้างตอนปิด (เลือกใช้ได้)
Environment=OPENCV_FFMPEG_CAPTURE_OPTIONS=rtsp_transport;tcp|max_delay;500000|stimeout;2000000

Restart=always
RestartSec=2

# รอปิดสั้นลง (ค่อยลดเหลือ 3 ถ้าทดสอบผ่านสม่ำเสมอ)
TimeoutStopSec=5
KillSignal=SIGTERM
KillMode=mixed
SendSIGKILL=yes

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
